# meta-openeuler/recipes-core/os-base/os-base/rc.local
#!/bin/bash

#/sbin/auditd
. /etc/init.d/auditd start
/sbin/syslogd
